# deploy/systemd/coupette-availability.service
[Unit]
Description=Coupette availability checker
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/home/deploy/coupette
EnvironmentFile=/home/deploy/coupette/.image-tag
EnvironmentFile=/etc/push-monitor/coupette-availability.env
ExecStart=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.prod.yml run --rm scraper python -m scraper availability
ExecStartPost=-/home/deploy/infra/scripts/push-monitor.sh ${PUSH_URL}
