# deploy/saq-scraper.service
[Unit]
Description=SAQ Sommelier scraper
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/saq-sommelier
EnvironmentFile=/opt/saq-sommelier/.image-tag
ExecStart=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.prod.yml run --rm scraper python -m scraper scrape
ExecStart=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.prod.yml run --rm scraper python -m scraper enrich
ExecStart=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.prod.yml run --rm scraper python -m scraper embed
