# deploy/saq-watches.timer
[Unit]
Description=Check watched SKU availability daily

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true

[Install]
WantedBy=timers.target
